# src/puppet/cassandra_pfpt/templates/cassandra-full-backup.service.erb
# /etc/systemd/system/cassandra-full-backup.service
# Managed by Puppet

[Unit]
Description=Cassandra Node Full Backup Service
Wants=cassandra.service
After=cassandra.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=<%= @full_backup_script_path %>