# src/puppet/cassandra_pfpt/templates/cassandra-full-backup.timer.erb
# /etc/systemd/system/cassandra-full-backup.timer
# Managed by Puppet

[Unit]
Description=Timer to schedule Cassandra node full backups

[Timer]
OnCalendar=<%= @full_backup_schedule %>
Persistent=true
Unit=cassandra-full-backup.service

[Install]
WantedBy=timers.target